# root/etc/s6-overlay/s6-rc.d/pproxy/finish
#!/usr/bin/with-contenv bash

pkill pproxy