# root/etc/s6-overlay/s6-rc.d/pproxy/finish
#!/usr/bin/with-contenv bash

ps -eo pid,args | grep /usr/bin/pproxy | \
    grep -v grep | awk '{print $1}' | xargs -r kill -9